# PKR/csye6225.service
[Unit]
Description=CSYE6225 Web Application
After=network.target

[Service]
User=csye6225
Group=csye6225
Type=simple
EnvironmentFile=/etc/csye6225/application-env
ExecStart=/usr/bin/java -jar /opt/csye6225/webapp.jar
Restart=always
RestartSec=10
WorkingDirectory=/opt/csye6225

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadOnlyDirectories=/
ReadWriteDirectories=/opt/csye6225 /var/log/webapp
PrivateDevices=true

[Install]
WantedBy=multi-user.target